# deploy/cloud/backend.service
[Unit]
Description=MiraLink FastAPI Backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=miralink
Group=miralink
WorkingDirectory=/opt/MiraLink/backend
EnvironmentFile=-/opt/MiraLink/deploy.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/MiraLink/.venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 8100
Restart=always
RestartSec=5
StandardOutput=append:/opt/MiraLink/logs/backend.log
StandardError=append:/opt/MiraLink/logs/backend.log

[Install]
WantedBy=multi-user.target
