# deploy/cloud/backend.service
[Unit]
Description=HerUnity FastAPI Backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=herunity
Group=herunity
WorkingDirectory=/opt/herunity/backend
EnvironmentFile=-/opt/herunity/deploy.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/herunity/.venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 8100
Restart=always
RestartSec=5
StandardOutput=append:/opt/herunity/logs/backend.log
StandardError=append:/opt/herunity/logs/backend.log

[Install]
WantedBy=multi-user.target
