# deploy/cloud/signalling.service
[Unit]
Description=MiraLink Python RenderStreaming Signalling Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=miralink
Group=miralink
WorkingDirectory=/opt/MiraLink/signalling
EnvironmentFile=-/opt/MiraLink/deploy.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/MiraLink/.venv/bin/python -u /opt/MiraLink/signalling/server.py
Restart=always
RestartSec=5
StandardOutput=append:/opt/MiraLink/logs/signalling.log
StandardError=append:/opt/MiraLink/logs/signalling.log

[Install]
WantedBy=multi-user.target
