# deploy/cloud/unity.service
[Unit]
Description=HerUnity Unity Linux Player
After=network-online.target herunity-backend.service herunity-signalling.service xvfb.service
Wants=network-online.target
Requires=herunity-backend.service herunity-signalling.service

[Service]
Type=simple
User=herunity
Group=herunity
WorkingDirectory=/opt/herunity/build
EnvironmentFile=-/opt/herunity/deploy.env
Environment=DISPLAY=:99
Environment=XDG_RUNTIME_DIR=/tmp
Environment=LD_LIBRARY_PATH=/opt/herunity/build:/opt/herunity/nvidia_libs
Environment=VK_ICD_FILENAMES=/etc/vulkan/icd.d/nvidia_icd.json
ExecStart=/opt/herunity/build/HerUnity.x86_64 -batchmode -RenderOffscreen -logfile /opt/herunity/logs/unity.log
Restart=always
RestartSec=8
StandardOutput=append:/opt/herunity/logs/unity.stdout.log
StandardError=append:/opt/herunity/logs/unity.stderr.log

[Install]
WantedBy=multi-user.target
