# scripts/unity_daemon.conf
; HerUnity Unity Pixel Streaming daemon
; Deploy: cp to /etc/supervisor/conf.d/unity_daemon.conf, then supervisorctl reread && supervisorctl update

[program:unity-pixel-streaming]
command=/data/HerUnity-Build/HerUnity.x86_64 -RenderOffscreen -PixelStreamingURL ws://127.0.0.1:8888 -logfile /data/logs/unity.log
directory=/data/HerUnity-Build
user=root
autostart=true
autorestart=true
startretries=5
startsecs=15
stopwaitsecs=10
stopsignal=TERM
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=/data/logs/unity_supervisor.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
environment=DISPLAY=:99,LD_LIBRARY_PATH="/data/HerUnity-Build:/data/nvidia_libs",XDG_RUNTIME_DIR="/tmp",HOME="/root"

; Backend daemon
[program:herunity-backend]
command=/data/backend/venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 8100
directory=/data/backend
user=root
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=5
stopsignal=TERM
redirect_stderr=true
stdout_logfile=/data/logs/backend_supervisor.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=2
environment=TTS_BACKEND="mock",AGENT_BACKEND="mock",ASR_BACKEND="mock",SERVER_PORT="8100"

; Signalling server daemon
[program:herunity-signalling]
command=node /data/signalling/server.js
directory=/data/signalling
user=root
autostart=true
autorestart=true
startsecs=3
stopwaitsecs=3
stopsignal=TERM
redirect_stderr=true
stdout_logfile=/data/logs/signalling_supervisor.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=1
environment=HTTP_PORT="7860",WS_PORT="8888"

; Xvfb daemon
[program:herunity-xvfb]
command=Xvfb :99 -screen 0 1920x1080x24 -ac +extension GLX
user=root
autostart=true
autorestart=true
startsecs=2
stopwaitsecs=2
redirect_stderr=true
stdout_logfile=/data/logs/xvfb_supervisor.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=1
